# Configs/.config/systemd/user/vyle-config.service
# ~/.config/systemd/user/vyle-config.service
[Unit]
Description=Vyle Project: TOML configuration watchdog
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/vyle-config
WorkingDirectory=%h
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
